# raspberry-pi-4-2gb-freebsd-router/tree/etc/rc.d/var-restore
#!/bin/sh

# REQUIRE: var
# BEFORE: FILESYSTEMS

. /etc/rc.subr

name="var-restore"
desc="Restore /var from backup"
stop_cmd=':'

load_rc_config $name

if [ -d /etc/var.backup ]; then
  cp -af /etc/var.backup/* /var/
fi
