# archlinux-vm-vbox-x86_64/tree/etc/xdg/autostart/lxqt-notifications-OPENBOX.desktop
[Desktop Entry]
Type=Application
Name=Notification Daemon
TryExec=lxqt-notificationd
Exec=lxqt-notificationd
OnlyShowIn=OPENBOX;
X-LXQt-Module=true
