# archlinux-vm-vbox-x86_64/tree/etc/systemd/system/sshdgenkeys.service
[Unit]
Description=SSH Key Generation

[Service]
ExecStart=/usr/bin/true
Type=oneshot
RemainAfterExit=yes
